# ip-ua.service
[Unit]
Requires=multi-user.target

[Service]
ExecStart=/path/to/Service
Restart=always